# conf.d/dropbear.conf
[program:dropbear]
command=service dropbear start
autorestart=false
