# conf.d/proxy.conf
[program:proxy]
command=/conf.d/proxy.sh
autorestart=false
