# conf.d/websockify.conf
[program:websockify]
command=/conf.d/websockify.sh
autorestart=true
