# conf.d/xterm.conf
[program:konsole]
command=konsole
autorestart=true
